# alpine-cron/crontab
# min	hour	day	month	weekday	command
*	*	*	*	*	touch /var/log/job_done.log
